# rootfs/etc/s6-overlay/s6-rc.d/init-usermod/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

UID=${UID:-911}
GID=${GID:-911}

echo "User Id $UID, Group id $GID"
usermod -u "$UID" nginx && groupmod -g "$GID" nginx
