# rootfs/etc/s6-overlay/s6-rc.d/svc-php82-fpm/run
#!/usr/bin/with-contenv bash

/usr/sbin/php-fpm82 --nodaemonize --fpm-config /etc/php82/php-fpm.conf