# rootfs/etc/s6-overlay/s6-rc.d/svc-stdout/run
#!/usr/bin/with-contenv bash

dockerize -stdout /var/www/storage/logs/laravel.log -stdout /var/log/nginx/error.log -stdout /var/log/php82/error.log -poll